# recipes-kolibre/kolibre-sample-client/files/media-mmc1.mount
[Unit]
Description=Mount mmc fat16 partition from mmc

[Mount]
What=/dev/mmcblk0p1
Where=/media/mmc1
Options=rw,errors=remount-ro
